# config/app.conf.example
[program:app]

command=/bin/bash -lc "bash /app/current/scripts/boot_production.bash"
